# templates/default/zookeeper.service.erb
[Unit]
Description = Zookeeper centralized coordination service
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.kzookeeper.user %>

WorkingDirectory=<%= node.zookeeper.base_dir %>
PIDFile = /var/run/zookeeper.pid
#ExecStart = <%= node.zookeeper.base_dir %>/bin/zookeeper-start.sh
ExecStart = ./bin/zkServer.sh start
#ExecStop = <%= node.zookeeper.base_dir %>/bin/zookeeper-stop.sh
ExecStop = ./bin/zkServer.sh stop

Type = forking

[Install]
WantedBy = multi-user.target
